# launchd/com.user.stockweekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.stockweekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd {{PROJECT_ROOT}} &amp;&amp; uv run scripts/weekly_loop.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/launchd_weekly_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/launchd_weekly_stderr.log</string>
</dict>
</plist>
